# Deployment/nginx/gunicorn.socket
[Unit]
Description=gunicorn socket

[socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target